# etc/systemd/var-cache.mount
[Unit]
DefaultDependencies=no
Before=umount.target
after=var-log.mount
Requires=var-log.mount
