# etc/systemd/var-log.mount
[Unit]
DefaultDependencies=no
Before=umount.target
after=systemd-journald.service
Requires=systemd-journald.service
